# crontab/route53-cron
# 占位文件 — 实际由 entrypoint.sh 动态生成
# 时区: Asia/Shanghai
*/1 * * * * /bin/bash -c 'source /app/state/container.env && /app/scripts/update-aws-route53-daemon.sh >> /app/logs/route53.log 2>&1'